# modules/redelinux-net/files/etc/profile.d/nfs_path.sh
#!/bin/sh
export PATH=$PATH:/global/bin
